# openrc/aria2cd
#!/sbin/openrc-run

name="aria2c"
description="Download utility that supports HTTP(S), FTP, BitTorrent, and Metalink"

command="/usr/bin/aria2c"
command_args="--conf-path=/etc/aria2/aria2.conf \
--dir=/home/young/Downloads \
--save-session=/home/young/.config/aria2/aria2.session"
pidfile="/run/${RC_SVCNAME}.pid"
command_user="young"

command_background=true

depend() {
    need net
    after net
}
